# mc-1.19.2.service
[Unit]
Description=A Paper Based 1.19.2 Minecraft Server

[Service]
ExecStart=~/paper-1.19.2/start-1.19.2.sh